# systemd/cerberus-celerybeat.service
[Unit]
Description=Cerberus Pi — Celery beat scheduler (nightly IPFS archival)
After=redis-server.service postgresql.service cerberus-celery.service
Wants=redis-server.service

[Service]
Type=simple
User=cerberus
Group=cerberus
WorkingDirectory=/opt/cerberus/backend
EnvironmentFile=/opt/cerberus/secrets/.env
ExecStart=/opt/cerberus/venv/bin/celery -A cerberus beat --loglevel=info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler
Restart=always
RestartSec=5
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
